# AC_Shield.Daemon/AC_Shield.service
[Unit]
Description=AC_Shield

[Service]
Type=notify
ExecStart=/usr/sbin/AC_Shield/AC_Shield.Daemon.exe

[Install]
WantedBy=multi-user.target